# deploy/busyq.service
[Unit]
Description=BusyQ - Real-time Crowd Heatmap App
Documentation=https://github.com/caioluchetti/busyq
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/busyq
Environment=PORT=8081
Environment=NODE_ENV=production
ExecStart=/usr/bin/node /var/www/busyq/server.js
Restart=on-failure
RestartSec=5
StandardOutput=append:/var/log/busyq.log
StandardError=append:/var/log/busyq.error.log

NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/www/busyq/.busyq-data

[Install]
WantedBy=multi-user.target
